# includedir/lxqt/etc/xdg/autostart/Compositor.desktop
[Desktop Entry]
Exec=xcompmgr -cfC
Name=Compositor
Type=Application
Version=1.0
