# includedir/lxqt/etc/xdg/autostart/Enable Super key shortcut.desktop
[Desktop Entry]
Exec=ksuperkey
Name=Enable Super key shortcut
Type=Application
Version=1.0
