# includedir/lxqt/etc/xdg/autostart/picom.desktop
[Desktop Entry]
Exec=picom -b
Name=Compositor
Type=Application
Version=1.0
