# includedir/lxqt/etc/xdg/autostart/set-lock-wallpaper.desktop
[Desktop Entry]
Exec=/usr/local/bin/set-lock-wallpaper
Name=Lock Wallpaper
Type=Application
Version=1.0
Name[en_US]=set-lock-wallpaper.desktop
