# includedir/lxqt/etc/xdg/autostart/xcompmgr.desktop
[Desktop Entry]
Exec=xcompmgr -cfC
Name=Compositor
Type=Application
Version=1.0
Name[en_US]=xcompmgr.desktop
