# includedir/lxqt/etc/skel/.config/autostart/enable-connman-gtk-tray.desktop
[Desktop Entry]
Exec=$HOME/.enable-connman-gtk-tray
Name=Enable Connman GTK tray icon
Type=Application
Version=1.0
